# ADMIN_SCRIPTS/MANAGER/rc.local
#!/bin/bash
 exit 0
 iptables-restore < /etc/iptables.sav
